# scripts/net.earthianlabs.cri-refresh-prices.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>net.earthianlabs.cri-refresh-prices</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/m3untold/Code/EarthianLabs/community-resilience-index/app/scripts/refresh-prices.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/m3untold/Code/EarthianLabs/community-resilience-index/app</string>

    <!-- Daily at 6am + 2pm — runs missed jobs when machine wakes -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key>
            <integer>6</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>14</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/m3untold/Code/EarthianLabs/community-resilience-index/app/scripts/logs/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/m3untold/Code/EarthianLabs/community-resilience-index/app/scripts/logs/launchd-stderr.log</string>

    <!-- Ensure PATH includes node -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
